# scripts/systemd/research-flow-frontend.service
[Unit]
Description=Research Flow Frontend (Next.js)
After=network.target research-flow-backend.service
Wants=research-flow-backend.service

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_GROUP
WorkingDirectory=/srv/research-flow/frontend
Environment="NODE_ENV=production"
Environment="NEXT_PUBLIC_API_BASE_URL=http://localhost:8000"
Environment="PORT=3000"
ExecStart=/usr/bin/npm run start -- --port 3000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=research-flow-frontend

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/research-flow

[Install]
WantedBy=multi-user.target

